# lesson6/roles/node_app/templates/systemd/node_app.service.j2
[Unit]
Description=Node.js Hello World App

[Service]
User={{ app_user }}
ExecStart=/usr/bin/node /opt/node_app/app.js
WorkingDirectory=/opt/node_app
Restart=always

[Install]
WantedBy=multi-user.target
